# systemd/qqmusic-first-sign.timer
[Unit]
Description=Run QQ Music medal check-in every day

[Timer]
OnCalendar=*-*-* 00:05:00
AccuracySec=1s
RandomizedDelaySec=0
Unit=qqmusic-first-sign.service

[Install]
WantedBy=timers.target
