# systemd/vae-checkin-fallback.timer
[Unit]
Description=Run Vae+ fallback check-in if primary did not leave a result

[Timer]
OnCalendar=*-*-* 00:00:45
OnCalendar=*-*-* 00:02:00
AccuracySec=1s
RandomizedDelaySec=0
Unit=vae-checkin-fallback.service

[Install]
WantedBy=timers.target
